# supervisor.dhcpd.conf
[program:dhcp-server]
command=/usr/sbin/dhcpd -cf /etc/dhcp/dhcpd.conf

